# docker/nginx/sscfg-ui/tests/e2e/misc/sscfg.ini
[program:dockerd]
command=/usr/local/bin/dockerd-entrypoint.sh

[program:sscfg]
command=/usr/local/bin/docker compose up
directory=/srv
environment=DOCKER_HOST=tcp://127.0.0.1:2376,DOCKER_CERT_PATH=/certs/client,DOCKER_TLS_VERIFY=1
redirect_stderr=true
stdout_logfile=/var/log/sscfg.log
